# archive.service
[Unit]
Description=Archive

[Service]
Type=oneshot
ExecStart=%h/.local/bin/archive.sh
WorkingDirectory=%h/.local/bin/
